# docker-mariadb/supervisor.conf.d/mariadb.conf
[program:mariadb]
command=/opt/replication-entrypoint.sh mysqld
autostart=true
